# config/includes.chroot_after_packages/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/launch.sh
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
